# deploy/spacegame-host.service
[Unit]
Description=spacegame authoritative sector host (pinned cells for the live galaxy)
# The game's authoritative sim runs INSIDE this process and talks to the mesh through the relay's
# own CE node, so it must come up after the node (ce-relay) and the hub are live.
After=network.target ce-relay.service ce-hub.service
Wants=ce-relay.service

[Service]
Type=simple
WorkingDirectory=/opt/ce-build/spacegame
# Spacegame is a first-class CE app, NOT a demo: it needs an always-on, server-class authoritative
# host the same way drift does. A browser-elected host renders but cannot run the authoritative sim
# (the wasm bundle ships only the renderer + local prediction), so without this the pinned sectors
# never advance. This pins the relay as the authoritative host for a plus-shaped block of sectors
# around the origin; every browser/native client renders and predicts against the state it publishes.
#
#   - The local CE node (ce-relay, `ce start --no-mine`) is the mesh transport on 127.0.0.1:8844.
#   - `host --ruleset live.json` watches the ruleset file: editing weapons / tech / tunables / SHADERS
#     in that file hot-reloads the whole live galaxy with no restart and no dropped players. The file
#     is seeded from the built-in ruleset by deploy/deploy.sh (`spacegame ruleset init`).
#   - Add or remove `--sector X_Y` to grow/shrink the pinned region; neighbours can also be placed on
#     other mesh nodes with `spacegame place` once more hosts join.
ExecStart=/opt/ce-build/spacegame/spacegame host \
  --sector 0_0 --sector 1_0 --sector -1_0 --sector 0_1 --sector 0_-1 \
  --hz 20 --ruleset /opt/ce-build/spacegame/live.json
Restart=always
RestartSec=3
# The host holds no secrets and only reaches the local node. Keep the box locked down; it only needs
# to read its own working dir (the binary + the live ruleset file the watcher reads).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ce-build/spacegame
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
